# scripts/deploy/supervisor-client.conf
[program:mattchat-client]
command = python /home/matt/www/mattbreen.net/breenchat/client/chat.py
directory = /home/matt/www/mattbreen.net/breenchat/client
user = matt
stdout_logfile = /home/matt/www/mattbreen.net/breenchat-client.log
redirect_stderr = true
environment=PATH="/home/matt/www/mattbreen.net/bin",LANG=en_US.UTF-8,LC_ALL=en_US.UTF-8

